# profile/airootfs/etc/systemd/system/harbor-compose.service
[Unit]
Description=Start Docker Compose stacks from NFS share
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/harbor-compose-up.sh
ExecStop=/usr/local/bin/harbor-compose-down.sh

[Install]
WantedBy=multi-user.target
