# profile/airootfs/etc/systemd/system/harbor-runner.service
[Unit]
Description=GitHub Actions Docker runner (harbor-srv-docker)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=exec
ExecStart=docker compose -f /etc/harbor-runner/compose.yaml up --remove-orphans
ExecStop=docker compose -f /etc/harbor-runner/compose.yaml down
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
