# profile/airootfs/etc/systemd/system/rpc-gssd.service.d/keytab-init.conf
[Service]
# Obtain a TGT from the KDC using the machine keytab after gssd starts.
# The NFS mount operation (sec=krb5i) requires root's credential cache
# (/tmp/krb5cc_0) even when rpc.gssd -n is used for per-UID file access.
ExecStartPost=/usr/bin/kinit -k nfs/harbor-srv@HARBOR.LOCAL
