# profile/airootfs/etc/systemd/system/rpc-gssd.service.d/krb5-ordering.conf
[Unit]
After=krb5-kdc.service
Requires=krb5-kdc.service
