# profile/airootfs/etc/systemd/system/rpc-gssd.service.d/krb5-ordering.conf
# rpc-gssd handles GSS/Kerberos token negotiation for NFS mounts. It must not
# start before the local KDC is ready, or GSS auth will fail on boot.
# Requires= (not just After=) makes rpc-gssd fail hard if the KDC doesn't
# start, rather than starting anyway and silently breaking Kerberos mounts.
[Unit]
After=krb5-kdc.service
Requires=krb5-kdc.service
