# profile/airootfs/etc/systemd/system/rpc-gssd.service.d/no-user-creds.conf
[Service]
# Use machine credentials for all UIDs, not per-user credential caches.
# Without this, non-root processes (e.g. compose-deploy) cannot establish
# a GSS context for sec=krb5i NFS access.
ExecStart=
ExecStart=/usr/sbin/rpc.gssd -n
