# profile/airootfs/etc/systemd/system/docker.service.d/nfs-dependency.conf
[Unit]
After=mnt-synology-harbor_srv-retry.service
Requires=mnt-synology-harbor_srv-retry.service
BindsTo=mnt-synology-harbor_srv.mount
